# techbleat-devops/backend/fastapi.service
# -----------------------------------------------------------------------------
# Techbleat Superstore - FastAPI systemd service
# Deploy to: /etc/systemd/system/fastapi.service
# -----------------------------------------------------------------------------
# Assignment: Deployment File Layout - /home/ec2-user/app/ = FastAPI backend
# Usage:
#   sudo systemctl daemon-reload
#   sudo systemctl start fastapi
#   sudo systemctl enable fastapi
# -----------------------------------------------------------------------------

[Unit]
Description=Techbleat FastAPI Application
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/app
EnvironmentFile=/home/ec2-user/app/.env
ExecStart=/home/ec2-user/app/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
